# src/canvos/overlay/files/etc/systemd/system/stylus-agent.service.d/bcm-sync.conf
[Service]
ExecStartPre=/usr/local/bin/bcm-sync-userdata.sh
