# src/canvos/overlay/files/etc/systemd/system/stylus-agent.service.d/bcm-sync.conf
[Unit]
# Wait for BCM cmd daemon to set hostname before syncing to Palette
After=cmd.service

[Service]
ExecStartPre=/usr/local/bin/bcm-sync-userdata.sh
